# goserver.conf
[program:goserver]
command=/home/vagrant/gosearch-master/server
autostart=true
autorestart=true
numprocs=1
